# automations/launchd/com.neo.codex_chronicle_sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.neo.codex_chronicle_sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/neo/.openclaw/workspace/.venv-main-safe/bin/python</string>
      <string>/Users/neo/.openclaw/workspace/scripts/sync_codex_chronicle.py</string>
      <string>--initial-tail</string>
      <string>40</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/neo/.openclaw/workspace</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/neo/.openclaw/logs/codex_chronicle_sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/neo/.openclaw/logs/codex_chronicle_sync.log</string>
  </dict>
</plist>
