# automations/launchd/com.neo.accountability_sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.neo.accountability_sweep</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/neo/.openclaw/workspace/.venv-main-safe/bin/python</string>
      <string>/Users/neo/.openclaw/workspace/scripts/accountability_sweep.py</string>
      <string>--api-url</string>
      <string>https://aiclone-production-32dc.up.railway.app</string>
      <string>--ready-age-minutes</string>
      <string>90</string>
      <string>--review-age-hours</string>
      <string>24</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/neo/.openclaw/workspace</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/neo/.openclaw/logs/accountability_sweep.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/neo/.openclaw/logs/accountability_sweep.log</string>
  </dict>
</plist>
